# crontab
0 5 * * * /usr/local/bin/node /app/src/bahn_price_scraper.js >> /var/log/cron.log 2>&1
# An empty line is required below this one for cron to work

